# aws_utilities/efs.automount
[Unit]
Description=AWS EFS automount

[Automount]
Where=/efs

[Install]
WantedBy=multi-user.target

